# server/namsbokasafn-efni.service
[Unit]
Description=Námsbókasafn Translation Pipeline Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/namsbokasafn-efni/server
Environment="NODE_ENV=production"
EnvironmentFile=/opt/namsbokasafn-efni/server/.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=namsbokasafn-efni

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/namsbokasafn-efni

[Install]
WantedBy=multi-user.target
